# limb-docs-indexer/debian/limb-docs-indexer.cron.d
5 4 * * sun root indexer --all --rotate -c /etc/limb-docs-indexer/sphinx.conf
